# KAIRO-FW/src/kairo-daemon/debian/postinst
#!/bin/sh
set -e

case "$1" in
    configure)
        # Create default config files if they don't exist
        mkdir -p /etc/kairo-fw
        [ -f /etc/kairo-fw/ai_hosts.txt ] || touch /etc/kairo-fw/ai_hosts.txt
        [ -f /etc/kairo-fw/blacklist.txt ] || touch /etc/kairo-fw/blacklist.txt

        # Enable and start the daemon
        systemctl daemon-reload
        systemctl enable kairo-fw.service
        systemctl restart kairo-fw.service
    ;;
esac

exit 0
