# KAIRO-FW/src/kairo-daemon/kairo-fw.service
[Unit]
Description=KAIRO Firewall Daemon (Ethics-free)
After=network.target

[Service]
ExecStart=/usr/bin/kairo-fw-daemon
Restart=always
User=root
Environment=LOG_LEVEL=info

[Install]
WantedBy=multi-user.target
